# dot.runit/sv/emacs/run
#!/usr/bin/env bash
exec emacs --fg-daemon
